# com.example.wxcollect.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  macOS launchd 定时任务模板。用法：
    1) 把下面 3 处 /PATH/TO/wechat-mp-collector 改成你的项目绝对路径
    2) 把 python3 路径改成你的（which python3 查看）
    3) 复制到 ~/Library/LaunchAgents/ 并 load（见 README「每天自动跑」）
  建议另存为 com.<你的名字>.wxcollect.plist。
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.example.wxcollect</string>

  <!-- 每天 09:30 跑一次。改 Hour/Minute 调整时间 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>30</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>/PATH/TO/wechat-mp-collector</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/PATH/TO/wechat-mp-collector/collect.py</string>
  </array>

  <!-- 日志 -->
  <key>StandardOutPath</key>
  <string>/PATH/TO/wechat-mp-collector/run.log</string>
  <key>StandardErrorPath</key>
  <string>/PATH/TO/wechat-mp-collector/run.log</string>

  <!-- 开机/休眠错过时补跑一次 -->
  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
